# client/systemd/network-monitor.service
[Unit]
Description=Network Monitor Client

[Service]
Type=simple
Restart=always
User=network-monitor
Group=network-monitor
WorkingDirectory=/usr/bin/network-monitor/client
Environment="ARGS=switch.lan,pi3.local,ping.projects.chrisjeakle.com"
ExecStart=/usr/bin/network-monitor/client/network-monitor $ARGS

[Install]
WantedBy=multi-user.target